# content/service/caddy/run
#!/bin/sh

exec 2>&1
exec cp /workdir/nginx.conf > /etc/nginx/conf.d/server.conf 
exec cat /etc/nginx/conf.d/server.conf 
exec rm -rf /etc/nginx/sites-enabled/default
exec nginx -g 'daemon off;'
